# scheduling/systemd/rsync-backup.service
[Unit]
Description=Backup files as configured in /etc/rsync-backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rsync-backup